# supervisord.conf
[supervisord]
nodaemon=true

[program:openfortivpn]
command=openfortivpn --config /etc/openfortivpn/config
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
